# paccache-clear/paccache-clear.service
[Unit]
Description=Clears the src/ cache for pacman and AUR packages

[Service]
Type=simple
ExecStart=/usr/local/bin/paccache-clear
User=root
